# clawdbot/scripts/services/clawdbot.service
[Unit]
Description=ClawdBot Telegram AI Trading Assistant
After=network.target
Wants=kalshi-edge-scanner.service yoshi-bridge.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/ClawdBot-V1
EnvironmentFile=-/root/ClawdBot-V1/.env
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=NODE_ENV=production

# Kill anything holding the port before start
ExecStartPre=-/usr/bin/fuser -k -9 18789/tcp
ExecStartPre=/bin/sleep 1

# Rebuild config from env vars before each start (picks up token changes)
# --merge preserves moltbot agent auth/bindings set by `moltbot agents add`
ExecStartPre=/usr/bin/python3 /root/ClawdBot-V1/scripts/rebuild-config.py --quiet --merge

# Ultimate-fix: Run diagnostics + auto-fix + calibration on start
# Uses reduced bars/forecasts for quick boot (full diagnostics via cron)
ExecStartPre=-/usr/bin/python3 -m scripts.forecaster.diagnose --bars 500 --forecasts 20 --auto-fix --output /root/ClawdBot-V1/data/diagnostics_report.json

# Kill port again right before start (ExecStartPre scripts take 15s+)
ExecStartPre=-/usr/bin/fuser -k -9 18789/tcp
ExecStartPre=/bin/sleep 1

# Start moltbot gateway
ExecStart=/usr/bin/moltbot gateway --port 18789

Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=clawdbot

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/root /home/root /tmp
ProtectHome=no

[Install]
WantedBy=multi-user.target
