# clawdbot/scripts/services/yoshi-bridge.service
[Unit]
Description=Yoshi Bridge — Outbox to Trading Core
After=network.target
PartOf=clawdbot.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/ClawdBot-V1
EnvironmentFile=-/root/ClawdBot-V1/.env
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=TRADING_CORE_URL=http://127.0.0.1:8000
Environment=YOSHI_OUTBOX_DIR=data/outbox

# Bridge forwards pending proposals from outbox to Trading Core /propose
ExecStart=/usr/bin/python3 scripts/yoshi-bridge.py --poll-interval 10 --max-send 100

Restart=always
RestartSec=15
StartLimitIntervalSec=600
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=yoshi-bridge

[Install]
WantedBy=multi-user.target
